# templates/default/iproute2.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:        iproute2
# Required-Start:    
# Required-Stop:     
# Default-Start:   2 3 4 5
# Default-Stop:    
# Short-Description: Start AGGREGATOR daemon
### END INIT INFO

# *****************************************************************
# Dynamic iproute2 configuration  --  DO NOT EDIT THIS FILE BY HAND
# Generated by chef - <%= @date %>
# *****************************************************************
PATH=/sbin:/bin:/usr/sbin:/usr/bin

. /lib/lsb/init-functions

ipRouteDeleteAll() {
 ip route| egrep -v "proto kernel|^[ ,	]"| while read r; do ip route del $r; done
}

do_start() {
 true
<% if @start %> ipRouteDeleteAll
<% @start.each do |i| %><%= " #{i}" %><% end %>
<%end%>
}

do_stop()  {
 true
<% if @stop %> ipRouteDeleteAll
<% @stop.each  do |i| %><%= " #{i}" %><% end %>
<%end%>
}

case "$1" in
start|restart|reload|force-reload)
 do_start || do_stop ;;
stop)
 do_stop ;;
status)
 ip route show ;;
*)
 echo "Syntaxe: $(basename $0) [start|stop]..." && false
esac

